# arch/nfs/media-nfsro.automount
# /etc/systemd/system/media-nfsro.automount
[Unit]
Description=Automount nfs Read Only

[Automount]
Where=/media/nfsro

[Install]
WantedBy=multi-user.target
